# packages/generate_index.timer
[Unit]
Description=Run generate_index.service everyday at 5:00

[Timer]
OnCalendar=*-*-* 5:00:00
Persistent=true

[Install]
WantedBy=timers.target
